# container.conf
description "has monit start stop container"
start on stopped rc RUNLEVEL=[2345]
stop on runlevel[!2345]
respawn
exec start-stop-daemon --start --exec /usr/local/bin/start-container start container